# code/supervisord.conf
[supervisord]
nodaemon=true

[program:track]
command=python /code/track.py

[program:server]
directory=/code
command=flask run --host=0.0.0.0
